# modules/bash/core/dist/.bash.d/profile.d/config
#!/bin/bash

shopt -s histappend
shopt -s extglob
shopt -s extquote
shopt -s cdspell
shopt -s dirspell
shopt -s globstar
shopt -s checkjobs
shopt -s checkwinsize

export PS1='\u@\H [\w]\n\$ '
export LANG=en_US.UTF-8
export LC_ALL=en_US.UTF-8
export HISTFILE=$HOME/.bash_history
export HISTFILESIZE=10000
export HISTSIZE=7500
export HISTCONTROL=ignoredups:erasedups
export HISTIGNORE=" *"
