# modules/bash/xtra/editor/dist/.bash.d/profile.d/editor
#!/bin/bash

export VISUAL=$HOME/.bash.d/bin/dx15_editor
export EDITOR=$HOME/.bash.d/bin/dx15_editor
