# ansible/roles/bash/templates/profile.d/editor
#!/bin/bash

export VISUAL="{{ root }}/.bash.d/bin/dx15_editor"
export EDITOR="{{ root }}/.bash.d/bin/dx15_editor"
