# bundles/dsouza@work/autorandr/dist/.bash.d/profile.d/autorandr_order.sh
#!/bin/sh

export AUTORANDR_ORDER="eDP-1 HDMI-1"
